# systemd/client/rsca.service
[Unit]
Description=RSCA Monitoring Command Client
BindTo=network.target

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/usr/sbin/rsca --watchdog
Type=simple
User=root
Group=root
Restart=always
WatchdogSec=120
